# homeassistant/root/etc/periodic/daily/dropbox_uploader
#!/bin/sh
/usr/bin/python /data/config/python_scripts/dropbox.py
